# deploy/systemd/sentinel-poll.timer
[Unit]
Description=Sentinel poll (ingest + categorize + policy/bill alerts), 4x/day

[Timer]
# 07:45, 12:45, 17:45, 21:45 local time — the PSD2 ~4/day unattended budget.
# This is what turns "I made a transaction" into an alert; latency is up to the
# gap between polls (a few hours), not instant. Use /sync for an on-demand pull.
OnCalendar=*-*-* 07,12,17,21:45:00
Persistent=true
Unit=sentinel-cron@poll.service

[Install]
WantedBy=timers.target
