# deploy/systemd/sentinel-backup.timer
[Unit]
Description=Sentinel nightly backup (02:30)

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true

[Install]
WantedBy=timers.target
